# saltmaster/supervisor-services.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:salt-master]
command=/usr/bin/salt-master -l debug
autostart=true
autorestart=true
