# saltmaster/supervisor-services.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:salt-master]
command=/home/vagrant/salt/scripts/salt-master -l debug -d
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true

[program:redis]
command=/usr/sbin/redis-server
autostart=true
autorestart=true

[program:ipython]
command=/usr/bin/ipython notebook
autostart=true
autorestart=true
