# minion/supervisor-services.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:salt-minion]
command=/usr/bin/salt-minion
autostart=true
autorestart=true

[program:dnsmasq]
command=/usr/bin/dnsmasq
autostart=true
autorestart=true