# puppetmaster/supervisor-services.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:puppet-master]
command=/usr/bin/puppet master --no-daemonize -d 
autostart=true
autorestart=true
